# slime-system/com.openclaw.slime-health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.slime-health</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/ki/.openclaw/workspace/slime-system/slime-agent.js</string>
        <string>--health</string>
    </array>
    <key>StandardOutPath</key>
    <string>/tmp/slime-health.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/slime-health-error.log</string>
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
</dict>
</plist>
